# SubversionEdge/files/etc/supervisor.d/csvnautostarter.conf
[program:csvnautostarter]
#command=/opt/csvnautostarter/csvnautostarter.sh
command=ls
user=collabnet